# .dev-hpo/php/supervisord.conf
[unix_http_server]
file = /var/run/supervisord.sock
chmod = 0766

[supervisord]
logfile = /tmp/supervisord.log
logfile_maxbytes = 10MB
logfile_backups = 10
loglevel = info
pidfile = /tmp/supervisord.pid
nodaemon = true
user = root
childlogdir = /log/supervisor

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl = unix:///var/run/supervisor.sock ; use a unix:// URL  for a unix socket

[program:php-fpm]
command = /usr/local/sbin/php-fpm
autostart = %(ENV_START_FPM)s
autorestart = true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:cron]
command = /usr/sbin/crond -f -l 7
autostart = %(ENV_START_CRON)s
autorestart = true
user = root
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:sf-messenger]
command = php /srv/api/bin/console messenger:consume -f 1 -l 150 -m 240M -t 2400 async
autostart = %(ENV_START_MESSENGER)s
autorestart = true
user = www-data
environment = HOME="/srv/api"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
